# arancino-rpi-scripts/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

#wifi-live-or-reset

# enable gpio as output to control reset pin on mikroe bus ch#1
echo 5 > /sys/class/gpio/export
echo out > /sys/class/gpio/gpio5/direction
echo 1 > /sys/class/gpio/gpio5/value

# configure and enable can0 interface
/usr/sbin/ip link set can0 type can bitrate 125000 triple-sampling on
/sbin/ifconfig can0 up

exit 0
